# cfg/crontab
MAILTO=""

0 * * * * sh -l $HOME/bin/secrets-backup
