# cfg/crontab
MAILTO=""

* * * * * sh -l $HOME/bin/secrets-backup
