# cfg/crontab
MAILTO=""

0 * * * * ~/bin/secrets-backup
